# systemd/tonto.timer
[Unit]
Description=Tonto Service Timer

[Timer]
Persistent=true
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=multi-user.target
